# nas/config/systemd/immich-ml-wake-proxy.service
[Unit]
Description=Wake PC worker stack and proxy Immich ML requests
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=nas
Environment=LISTEN_HOST=0.0.0.0
Environment=LISTEN_PORT=3004
Environment=PC_ML_URL=http://pc:3003
Environment=LOCAL_ML_URL=http://127.0.0.1:3003
Environment=ENSURE_SCRIPT=/data/docker/compose/nightly-orchestrator/pc-worker-ensure.sh
Environment=NIGHT_ONLY=true
Environment=NIGHT_START=04:00
Environment=NIGHT_END=10:00
ExecStart=/usr/bin/python /data/docker/compose/nightly-orchestrator/immich-ml-wake-proxy.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
