# nas/config/systemd/nas-uptime-ledger.service
[Unit]
Description=Update NAS availability ledger
After=systemd-journald.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/nas-uptime-ledger
