# nas/config/systemd/nas-url-queue-notify.service
[Unit]
Description=Send URL queue completion notifications
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/nas-url-queue-notify
