# nas/config/systemd/nas-weekly-digest.service
[Unit]
Description=Generate and send the weekly NAS digest
RequiresMountsFor=/data /mnt/disk1 /mnt/disk2 /mnt/disk3 /mnt/parity
After=local-fs.target docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/nas-weekly-digest
