# nas/config/systemd/nas-parity-spindown.service
[Unit]
Description=Set the NAS parity disk standby timeout
RequiresMountsFor=/mnt/parity
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/nas-parity-spindown
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
