# nas/config/systemd/nas-weekly-digest.timer
[Unit]
Description=Send the weekly NAS digest

[Timer]
OnCalendar=Sun 09:00
RandomizedDelaySec=15m
Persistent=true
Unit=nas-weekly-digest.service

[Install]
WantedBy=timers.target
