# nas/config/systemd/nas-kernel-maintenance-reminder.service
[Unit]
Description=Send semiannual NAS kernel maintenance reminder
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/nas-kernel-maintenance-reminder
