# nas/config/systemd/snapraid-scrub.timer
[Unit]
Description=Run SnapRAID partial scrub weekly

[Timer]
# Reserve early Sunday for parity verification, before the PC-worker window.
OnCalendar=Sun 00:30
RandomizedDelaySec=15m
Persistent=true

[Install]
WantedBy=timers.target
