# nas/config/systemd/tdarr-wake-monitor.timer
[Unit]
Description=Poll Tdarr for queued work and idle PC worker cleanup

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=15s
Unit=tdarr-wake-monitor.service

[Install]
WantedBy=timers.target
