# nas/config/systemd/docker.service.d/wait-for-data.conf
# Generated by NAS bootstrap.
# nftables protects host services. Docker published ports are explicit exposure.
# Compose files should avoid broad 0.0.0.0 bindings unless intentionally reviewed.
# Prefer LAN-only, Tailscale-only, or localhost bindings for admin UIs.
[Unit]
RequiresMountsFor=/data
After=network-online.target
Wants=network-online.target
